# guia5/ej3_guia5.timer
[Unit]
Description=Timer ej3_guia5

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
