# guia5/ej5_guia5.timer
[Unit]
Description=Timer ej5_guia5

[Timer]
OnCalendar=*:0/2

[Install]
WantedBy=timers.target
